# roles/makelunch/templates/upstart/upstart.conf
description "{{service}}"

start on started networking
stop on runlevel [016]

limit nofile 1000000 1000000

respawn
respawn limit 5 30

setuid {{user}}
setgid {{user}}

chdir /home/{{user}}/{{service}}/build

env NODE_ENV={{node_env}}
env ROOT_URL={{root_url}}
env MONGO_URL=mongodb://localhost:27017/{{service}}
env PORT=3030

script
  METEOR_SETTINGS="$(cat /home/{{user}}/{{service}}/settings.json)"
  export METEOR_SETTINGS
  exec node main.js
end script
